# root/etc/systemd/system/guestfs-firstboot.service.d/firstboot-override.conf
[Unit]
After=nethserver-config-network.service
