# root/etc/cron.daily/apu-fstrim
#!/bin/bash

STATUSFILE=/var/tmp/last_trim
/bin/fstrim -v /boot >$STATUSFILE 2>&1
/sbin/fstrim -v / >>$STATUSFILE 2>&1
